# supervisord.conf
[supervisord]
nodaemon=true

[program:ddns]
command=/usr/src/app/ddns
autostart=true
startsecs=10
startretries=3
priority=1
redirect_stderr=true
stdout_logfile=/usr/src/app/ddns.log
